# old-stuff/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

killall nash-hotplug
if [ -f /etc/ec2.d/firstboot ]; then
  rm -rf /root/.ssh
  mkdir /root/.ssh
  chmod 700 /root/.ssh
  curl --retry 3 --retry-delay 5 --silent --fail -o /root/.ssh/authorized_keys http://169.254.169.254/1.0/meta-data/public-keys/0/openssh-key
  chmod 400 /root/.ssh/authorized_keys
  
  /usr/bin/ruby /etc/ec2.d/prepare.rb
  
  rm -f /etc/ec2.d/firstboot
fi
